# systemd/acts.service
[Unit]
Description=Acts Tarsnap backup

[Service]
Type=simple
ExecStart=/bin/acts
